# templates/default/ogawa.service.erb
[Unit]
Description=Ogawa Service
After=syslog.target

[Service]
Type=simple
User=<%= @user %>
WorkingDirectory=<%= @dir %>
ExecStart=/usr/bin/env python <%= @script %>
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
